# learn/ishan/iRaidersRobot2019/coprocessor/systemd/vision-watcher.path
[Path]
PathModified=/home/pi/vision

[Install]
WantedBy=multi-user.target
